# magicforce68.service
[Unit]
Description=MagicForce 68 Keyboard spacefn service

[Service]
Type=simple
ExecStart=/opt/spacefn/spacefn /opt/spacefn/mf68.cfg

